# wyoming-soniox/marc-wyoming-soniox.service
[Unit]
Description=MARC voice — Wyoming bridge for Soniox real-time STT
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.hermes/.env
WorkingDirectory=%h/projects/marc-voice/wyoming-soniox
ExecStart=%h/projects/marc-voice/wyoming-soniox/venv/bin/python -m wyoming_soniox --uri tcp://0.0.0.0:10300 --max-endpoint-delay-ms 1800 --terms-file %h/projects/marc-voice/wyoming-soniox/marc-context-terms.txt
Restart=on-failure
RestartSec=5
# Soniox key + small footprint; nothing else exotic.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
